# scripts/launchd/com.knut.clickup-init-week.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.knutkirkhorn.clickup-init-week</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd /Users/knutkirkhorn/dev/dotfiles &amp;&amp; set -a &amp;&amp; source .env &amp;&amp; set +a &amp;&amp; bun run scripts/init-clickup-week.ts; exit_code=$?; if [ $exit_code -ne 0 ]; then osascript -e 'display alert "ClickUp weekly init failed" message "Check /tmp/clickup-init-week.err.log" as critical buttons {"OK"} default button "OK"'; afplay /System/Library/Sounds/Basso.aiff; fi; exit $exit_code</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>5</integer>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/clickup-init-week.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/clickup-init-week.err.log</string>
  </dict>
</plist>
